# op/scripts/shadowsocks/shadowsocks.ini
[program:shadowsocks]
command=#ssserver# -c /etc/shadowsocks.json
autorestart=true
startsecs=2
user=shadowsocks
