# op/scripts/systemd.service
[Service]
ExecStart=/usr/bin/yes
User=oott123
Restart=always
RestartSec=5
Environment=A=B
Environment=C=D

[Install]
WantedBy=multi-user.target
